# initramfs_root/etc/inittab
::askfirst:-/bin/sh
::restart:/sbin/init
::shutdown:/bin/umount -a -r
